# recipes/set_motd/files/etc/update-motd.d/99-welcome
#!/bin/sh

echo "- Welcome to Pytagora'switches World! -"
