# recipes/set_crontab_for_user/templates/conf/pitanetes/crontab/crontab.erb
@hourly /usr/bin/docker ps -f name=coredns -f publish=53/udp -q | grep '[0-9|a-z]' || /usr/bin/docker run -d --rm -v /home/<%= @user %>/conf/pitanetes/coredns/:/root/ --name coredns -p 53:53 -p 53:53/udp coredns/coredns -conf /root/Corefile
